# .config/sv/dockerd-rootless
#!/sbin/openrc-run

description="Docker Daemon (Rootless)"

command="/usr/bin/dockerd-rootless"

supervisor=supervise-daemon
error_logger="logger -t '${RC_SVCNAME}' -p daemon.error"

# note: add to ~/.rc
# export DOCKER_HOST=unix://$XDG_RUNTIME_DIR/docker.sock
